# service/fastcgi-mono-server4/run
#!/bin/sh
exec 2>&1
exec /usr/bin/fastcgi-mono-server4 /applications=/:/var/www /socket=tcp:127.0.0.1:9000
